# .local/share/runit/sv/ssh-agent/run
#!/bin/sh

exec >/dev/null
exec ssh-agent -Da "${XDG_DATA_HOME:-$HOME/.local/share}"/ssh-agent.sock
