# home/bazz/.config/cron.weekly/pictures
#!/bin/sh
exec rsync -avu /home/bazz/Pictures/ /mnt/wd1tb/backup/pictures/
