# etc/rc.local
#!/bin/sh

# I don't want that mouse wakes up my computer
echo "disabled" >/sys/bus/usb/devices/1-10/power/wakeup
