# home/bazz/.config/service/ssh-agent/finish
#!/bin/sh

rm -f "${XDG_DATA_HOME:-$HOME/.local/share}"/ssh-agent.sock
exit 0
