# home/bazz/.config/service/ssh-agent/run
#!/bin/sh

exec >/dev/null
exec chpst -e "$TURNSTILE_ENV_DIR" ssh-agent -Da "${XDG_DATA_HOME:-$HOME/.local/share}"/ssh-agent.sock
