# .config/service/syncthing/run
#!/bin/sh

exec chpst -e "$TURNSTILE_ENV_DIR" syncthing serve --no-browser
